# nodejs.conf
[program:nodesprocess]
command=/usr/bin/node /app/api/server.js
directory=/
autostart=true
autorestart=true
startretries=3